# installer-fs/etc/inittab
::sysinit:/etc/init.d/rcS

tty1::once:/install
